# templates/default/rootfs/chef-install.conf.erb
start on net-device-up IFACE=public
exec /usr/sbin/install-chef.sh
